# english-analysis.service
[Unit]
Description=English Analysis Application
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/tmp/English-master_-english-analysis
Environment="NODE_ENV=production"
Environment="PORT=3001"
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=english-analysis

[Install]
WantedBy=multi-user.target

